# modules/pilights/templates/disable_network_led.service.erb
[Unit]
Description=Turn off network activity light
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/disable_network_led

[Install]
WantedBy=default.target
